# templates/systemd-beanstalk.service.erb
# systemd unit file for howalarming-beanstalk
# (generated via puppet-howalarming)
#

[Unit]
Description=Beanstalk queue for howalarming

[Service]
Type=simple
ExecStart=<%= @beanstalk_binary -%> -l 127.0.0.1 -p <%= @beanstalk_port %>
RestartSec=5
Restart=on-failure

[Install]
WantedBy=multi-user.target
